# bug-bounty-app.service
[Unit]
Description=Bug Bounty Application
After=network.target

[Service]
Type=simple
User=www-data
EnvironmentFile=/opt/bug-bounty-app/.env
WorkingDirectory=/opt/bug-bounty-app
ExecStart=/usr/bin/node backend/server.js
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
